# init/supervisord/home-automation.ini
[program:home-automation]
command=docker compose up
directory=/opt/home-automation
user=pi
group=pi
autostart=true
autorestart=true
startretries=3
startsecs=10
stopsignal=TERM
stopwaitsecs=30
killasgroup=true
stopasgroup=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/home-automation.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=HOME="/opt/home-automation",COMPOSE_PROJECT_NAME="home-automation"

[program:home-automation-prestart]
command=/bin/bash -c "while ! docker info >/dev/null 2>&1; do sleep 1; done; docker compose pull --quiet"
directory=/opt/home-automation
user=pi
group=pi
autostart=true
autorestart=false
startretries=1
startsecs=0
priority=100
stdout_logfile=/var/log/supervisor/home-automation-prestart.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
