# init/supervisord/home-automation-standalone.ini
[program:home-automation-standalone]
command=/opt/home-automation/bin/home-automation-server
directory=/opt/home-automation
user=pi
group=pi
autostart=true
autorestart=true
startretries=3
startsecs=5
stopsignal=TERM
stopwaitsecs=30
redirect_stderr=true
stdout_logfile=/var/log/supervisor/home-automation-standalone.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=10
environment=HOME="/opt/home-automation",LOG_LEVEL="info"
