# init/upstart/home-automation.conf
# Home Automation System - Upstart Job

description "Home Automation System with Docker Compose"
author "Home Automation Team"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [!2345]

# Restart the service if it dies
respawn
respawn limit 10 5

# Set working directory and user
chdir /opt/home-automation
setuid pi
setgid pi

# Environment setup
env HOME=/opt/home-automation
env COMPOSE_PROJECT_NAME=home-automation

# Load environment files
source /etc/default/home-automation || true
source /opt/home-automation/.env || true

# Pre-start script
pre-start script
    # Ensure Docker is running
    while ! docker info >/dev/null 2>&1; do
        sleep 1
    done
    
    # Pull latest images
    docker compose pull --quiet || true
end script

# Main script
script
    exec docker compose up
end script

# Post-stop script
post-stop script
    docker compose down || true
end script

# Console logging
console log
